# deploy/tgchannelSeb-web.service
[Unit]
Description=Aromat Day Web Admin Panel
After=network.target postgresql.service

[Service]
WorkingDirectory=/opt/tgchannelSeb
EnvironmentFile=/opt/tgchannelSeb/.env
ExecStart=/opt/tgchannelSeb/.venv/bin/python -m app.web
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
